# systemd/.config/systemd/user/dotfiles-backup.timer
[Unit]
Description=Hourly dotfiles backup while working

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=timers.target
